# root/etc/profile.d/10-desktop
#!/bin/bash

# Disable desktop wallpaper
xfconf-query  \
  --channel xfce4-desktop \
  --property /backdrop/screen0/monitorrdp0/workspace0/image-style \
  --set 0